# rpm/rtpproxy.socket
[Unit]
Description=RTPproxy controlling socket

[Socket]
ListenStream=/run/rtpproxy.sock
Accept=false

[Install]
WantedBy=sockets.target

